# xserver-vps-renew.service
[Unit]
Description=XServer VPS Auto Renew
After=network.target

[Service]
Type=oneshot
User=exedev
# 修改为你的实际路径
WorkingDirectory=/home/exedev/xserver-vps-renew
Environment=DISPLAY=:99
ExecStart=/bin/bash -c 'Xvfb :99 -screen 0 1280x900x24 & sleep 2 && /home/exedev/.local/bin/uv run python /home/exedev/xserver-vps-renew/renew.py; pkill -f "Xvfb :99"'
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
